# server-setup/supervisor.ini
; this is the configuration for supervisor
; "fab webserver deploy" will copy it to /etc/supervisor/project_name.ini
; check if there’s a line like the following in your supervisord.conf, [supervisord] section:
; environment=LANG=de_DE.UTF-8,LC_ALL=de_DE.UTF-8,LC_LANG=de_DE.UTF-8

[program:project_name]
command=/var/www/project_name/bin/gunicorn project_name.wsgi:application --bind=127.0.0.1:800x --workers=2 -c gunicorn-settings.py
; process_name=appserver
autorestart=true
user=project_name
stdout_logfile=/var/www/project_name/logs/supervisor.log
stdout_logfile_maxbytes=1MB
;stdout_logfile_backups=10
stderr_logfile=/var/www/project_name/logs/supervisor_error.log
stderr_logfile_maxbytes=1MB
;stderr_logfile_backups=10
directory=/var/www/project_name/releases/current
environment=HOME=/home/project_name
